# launchd/com.lana-bot.cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lana-bot.cycle</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/ada/lana-bot/scripts/cycle.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ada/lana-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ada</string>
    </dict>

    <!-- Every 30 minutes, matching strategy.toml cycle_minutes. -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ada/lana-bot/logs/cycle.launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ada/lana-bot/logs/cycle.launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
