# launchd/com.lana-bot.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lana-bot.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>scripts/monitor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ada/lana-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ada</string>
    </dict>

    <!-- Long-running stop-loss daemon — keep alive, but don't hammer restart. -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ada/lana-bot/logs/monitor.launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ada/lana-bot/logs/monitor.launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
